# crontabs/root
0 18 * * * /usr/local/bin/gsutil_wrapper
